# tests/nix_tests/data/evmos_9000-1/tasks.ini
[program:evmos_9000-1-node0]
autostart = true
autorestart = true
redirect_stderr = true
startsecs = 3
directory = %(here)s/node0
command = evmosd start --home . --trace
stdout_logfile = %(here)s/node0.log

[program:evmos_9000-1-node1]
autostart = true
autorestart = true
redirect_stderr = true
startsecs = 3
directory = %(here)s/node1
command = evmosd start --home . --trace
stdout_logfile = %(here)s/node1.log

