# tests/nix_tests/data/tasks.ini
[include]
files = %(here)s/cosmoshub-1/tasks.ini %(here)s/evmos_9000-1/tasks.ini

[supervisord]
pidfile = %(here)s/supervisord.pid
nodaemon = true
logfile = /dev/null
logfile_maxbytes = 0
strip_ansi = true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = %(here)s/supervisor.sock

[supervisorctl]
serverurl = unix://%(here)s/supervisor.sock

[program:relayer-demo]
autostart = false
autorestart = true
redirect_stderr = true
startsecs = 3
directory = %(here)s
command = hermes --config relayer.toml start
stdout_logfile = %(here)s/relayer-demo.log

