# tests/nix_tests/data/cosmoshub-1/tasks.ini
[program:cosmoshub-1-node0]
autostart = true
autorestart = true
redirect_stderr = true
startsecs = 3
directory = %(here)s/node0
command = gaiad start --home . 
stdout_logfile = %(here)s/node0.log

[program:cosmoshub-1-node1]
autostart = true
autorestart = true
redirect_stderr = true
startsecs = 3
directory = %(here)s/node1
command = gaiad start --home . 
stdout_logfile = %(here)s/node1.log

